# roles/apnscp/install-services/files/apnscp-housekeeping.timer
[Unit]
Description=Run apnscp housekeeping

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=60min

[Install]
WantedBy=timers.target