# roles/apnscp/install-services/files/rc-local.service
[Unit]
Description=/etc/rc.local compatibility
ConditionPathExists=/etc/rc.local

[Service]
Type=oneshot
ExecStart=/bin/sh /etc/rc.local
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
